# luxebox/daily_driver/launchd/com.luxebox.litellm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.luxebox.litellm</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uvx</string>
        <string>--from</string>
        <string>litellm[proxy]>=1.83,!=1.82.7,!=1.82.8</string>
        <string>litellm</string>
        <string>--config</string>
        <string>/Users/michaeltimpe/Downloads/luxebox/daily_driver/litellm_bridge.yml</string>
        <string>--port</string>
        <string>4000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/michaeltimpe/Downloads/luxebox</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>/Users/michaeltimpe/Library/Logs/luxebox/litellm.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/michaeltimpe/Library/Logs/luxebox/litellm.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
